# build/deb/etc/supervisor/conf.d/rukmini.conf
[program:rukmini]
directory=/var/lib/rukmini/src/rukmini
command=/var/lib/rukmini/src/rukmini/rukmini
pidfile = /var/run/rukmini/rukmini.pid
process_name=rukmini
autostart=true
autorestart=true
umask=777
priority=999
autostart=true
autorestart=true
startsecs=5
startretries=100
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=5
redirect_stderr=true
stdout_logfile=/var/log/rukmini/rukmini.out
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB