# build/deb/DEBIAN/postinst
#! /bin/bash -e
#Package Configuration
PACKAGE=rukmini
PACKAGE_USER=rukmini
PACKAGE_GROUP=rukmini
PF_UID=9075
PF_GID=9001

if [ "$1" == "configure" ] ; then
    #Creating group if it doesnt exist
    if ! getent group ${PACKAGE_GROUP} > /dev/null; then
        groupadd -g "${PF_GID}" ${PACKAGE_GROUP}
    fi
    #Creating user if it doesnt exist
    if ! getent passwd ${PF_UID} > /dev/null; then
        adduser --system --uid ${PF_UID} --home /var/lib/${PACKAGE} --no-create-home \
            --ingroup ${PACKAGE_GROUP} --disabled-password --shell /bin/false \
            ${PACKAGE_USER}
    fi
    if ! getent passwd ${PF_UID} > /dev/null; then
        adduser --system --uid ${PF_UID} --home /var/lib/${PACKAGE} --no-create-home \
            --ingroup ${PACKAGE_GROUP} --disabled-password --shell /bin/false \
            ${PACKAGE_USER}
    fi
    chown -R ${PF_UID}:${PF_GID} /var/lib/${PACKAGE}
    if [ -d "/var/log/${PACKAGE}" ]; then
 	    touch /var/log/${PACKAGE}/${PACKAGE}.log
 	else
 	    mkdir -p /var/log/${PACKAGE}
 	    touch /var/log/${PACKAGE}/${PACKAGE}.log
    fi

    chmod -R a+rwx /var/log/${PACKAGE}
    chown -R ${PF_UID}:${PF_GID} /var/log/${PACKAGE}
    chmod a+x /etc/init.d/${PACKAGE}
    chmod -R a+rwx /var/lib/${PACKAGE}/src

    # Compile the binary now (double meh.. tripple meh)
    export PATH=/usr/local/go/bin:$PATH
    export GOPATH=/var/lib/rukmini
    cd /var/lib/rukmini/src/rukmini
    echo "Compiling rukmini service..."
    go build
    echo "Compilation complete!"

    if [ ! -d "/usr/${PACKAGE}/bin" ]; then
        mkdir -p /usr/${PACKAGE}/bin
    fi
    #Compile scripts
    echo "Compiling scripts"
    cd /var/lib/rukmini/src/rukmini/scripts/
    echo "Compiling consumer script"
    go build consumer.go
    cp /var/lib/rukmini/src/rukmini/scripts/consumer /usr/rukmini/bin
    echo "Compiled consumer"
    echo "Compiling producer script"
    go build producer.go
    cp /var/lib/rukmini/src/rukmini/scripts/producer /usr/rukmini/bin
    echo "Compiled producer"
    echo "Completed scripts!!!"

    chmod -R a+rwx /var/lib/${PACKAGE}/src/${PACKAGE}/${PACKAGE}
    chmod 777 /etc/init.d/producer
    chmod 777 /etc/init.d/consumer

    sudo chkconfig --level 2345  ${PACKAGE} on
    sudo /etc/init.d/${PACKAGE} stop || true
    sleep 5
    if [ -x "/etc/init.d/${PACKAGE}" ]; then
	    sudo /etc/init.d/${PACKAGE} start || true
	    sleep 5
    fi
fi
exit 0