# build/deb/DEBIAN/postinst
#! /bin/bash -e
#Package Configuration
PACKAGE=rukmini
PACKAGE_USER=rukmini
PACKAGE_GROUP=rukmini
PF_UID=9075
PF_GID=9001

if [ "$1" == "configure" ] ; then
    #Creating group if it doesnt exist
    if ! getent group ${PACKAGE_GROUP} > /dev/null; then
        groupadd -g "${PF_GID}" ${PACKAGE_GROUP}
    fi
    #Creating user if it doesnt exist
    if ! getent passwd ${PF_UID} > /dev/null; then
        adduser --system --uid ${PF_UID} --home /var/lib/${PACKAGE} --no-create-home \
            --ingroup ${PACKAGE_GROUP} --disabled-password --shell /bin/false \
            ${PACKAGE_USER}
    fi
    if ! getent passwd ${PF_UID} > /dev/null; then
        adduser --system --uid ${PF_UID} --home /var/lib/${PACKAGE} --no-create-home \
            --ingroup ${PACKAGE_GROUP} --disabled-password --shell /bin/false \
            ${PACKAGE_USER}
    fi
    chown -R ${PF_UID}:${PF_GID} /var/lib/${PACKAGE}
    if [ -d "/var/log/${PACKAGE}" ]; then
 	    touch /var/log/${PACKAGE}/${PACKAGE}.log
 	    chmod -R 777 /var/log/${PACKAGE}/${PACKAGE}.log
 	else
 	    mkdir -p /var/log/${PACKAGE}
 	    touch /var/log/${PACKAGE}/${PACKAGE}.log
 	    chmod -R 777 /var/log/${PACKAGE}/${PACKAGE}.log
    fi
    chmod 777 /etc/init.d/${PACKAGE}
    chmod -R 777 /var/lib/${PACKAGE}/src
    if [ -d "/var/run/${PACKAGE}" ]; then
 	    chmod -R 777 /var/run/${PACKAGE}
 	else
 	    mkdir -p /var/run/${PACKAGE}
    fi
    touch /var/run/${PACKAGE}/${PACKAGE}.pid
    chmod -R 777 /var/run/${PACKAGE}
    chown -R ${PF_UID}:${PF_GID} /var/run/${PACKAGE}

    chmod 777 /etc/init.d/${PACKAGE}
    chmod -R 777 /var/lib/${PACKAGE}/src

    ## Delete job for package user. Must be run as root user ##
    echo "Adding cron job to clear temp..."
    crontab -r -u ${PACKAGE_USER} || true
    { echo '*/2 * * * * /etc/init.d/rukmini cleartemp'; } | crontab -u ${PACKAGE_USER} - || true
    crontab -u ${PACKAGE_USER} -l || true
    # Compile the binary now (double meh.. tripple meh)
    export PATH=/usr/local/go/bin:$PATH
    export GOPATH=/var/lib/rukmini
    cd /var/lib/rukmini/src/rukmini
    echo "Compiling rukmini service..."
    go build
    echo "Compilation complete!"

    chmod -R a+rwx /var/lib/${PACKAGE}/src/${PACKAGE}/${PACKAGE}

    sudo chkconfig --level 2345  ${PACKAGE} on
    sudo /etc/init.d/${PACKAGE} terminate || true

    cp -r /etc/default/haproxy-rukmini /etc/default/haproxy || true
    cat /dev/null >/var/log/${PACKAGE}/${PACKAGE}.log
    sleep 5
    if [ -x "/etc/init.d/${PACKAGE}" ]; then
	    sudo /etc/init.d/${PACKAGE} start || true
	    sleep 5
    fi

    echo "Renaming logrotate and rsyslog files"
    cp /etc/logrotate.d/rukmini-haproxy /etc/logrotate.d/haproxy
    cp /etc/rsyslog.d/rukmini-49-haproxy.conf /etc/rsyslog.d/49-haproxy.conf
fi
exit 0